# linux/home/.config/systemd/user/artifact-cloud-maintenance.timer
[Unit]
Description=Daily artifact-cloud maintenance

[Timer]
OnCalendar=*-*-* 03:17:00
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
