# linux/home/.config/systemd/user/artifact-cloud-drill.timer
[Unit]
Description=Monthly artifact-cloud restore drill

[Timer]
OnCalendar=*-*-01 04:17:00
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
