# linux/home/.config/systemd/user/artifact-cloud-maintenance.service
[Unit]
Description=Back up and verify the personal artifact cloud

[Service]
Type=oneshot
ExecStart=%h/.local/bin/artifact-cloud maintain
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/share/artifact-cloud %h/.local/state/artifact-cloud
